# etc/inittab
::sysinit:/etc/start
::respawn:cat /dev/hidraw0 > /dev/hidg0
